# templates/automatic-upgrades.timer.j2
{{ ansible_managed|comment }}

[Unit]
Description=automatic system upgrade timer

[Timer]
OnCalendar={{ timer_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
